# systemd/tascam-remap.service
[Unit]
Description=Create TASCAM Model 12 per-channel virtual sources

[Service]
Type=simple
ExecStart=%h/pulse/scripts/tascam-remap-watch.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
